# src/etc/services-config/supervisor/supervisord.d/haproxy-wrapper.conf
[program:haproxy-wrapper]
priority = 100
command = /usr/sbin/haproxy-wrapper
startsecs = 0
autorestart = true
redirect_stderr = true
stdout_logfile = /var/log/haproxy.log
stdout_events_enabled = true