# etc/services-config/supervisor/supervisord.d/haproxy-bootstrap.conf
[program:haproxy-bootstrap]
priority = 6
command = /usr/sbin/haproxy-bootstrap
startsecs = 0
startretries = 0
autorestart = false
redirect_stderr = true
stdout_logfile = /var/log/haproxy.log
stdout_events_enabled = true